# signal_scan.crontab
# LIFFE
49 17 * * 1-5 /cygdrive/c/tmp/backfill.rb
50 17 * * 1-5 /cygdrive/c/tmp/signal_scan.rb SS_RTL SS_RT_CTL

# EUREX
49 20 * * 1-5 /cygdrive/c/tmp/backfill.rb
50 20 * * 1-5 /cygdrive/c/tmp/signal_scan.rb EI_RTL EI_RT_CTL SB_RT_CTS

# GLOBEX
49 21 * * 1-5 /cygdrive/c/tmp/backfill.rb
50 21 * * 1-5 /cygdrive/c/tmp/signal_scan.rb ED_RT_CTL ED_ZS_L
