# signal_scan.crontab
# Signal scan
50 21 * * * 1-5 * /cygdrive/c/tmp/signal_scan.rb ED_RT_CTL
50 20 * * * 1-5 * /cygdrive/c/tmp/signal_scan.rb EI_RTL
50 20 * * * 1-5 * /cygdrive/c/tmp/signal_scan.rb EI_RT_CTL
50 17 * * * 1-5 * /cygdrive/c/tmp/signal_scan.rb SS_RTL
50 17 * * * 1-5 * /cygdrive/c/tmp/signal_scan.rb SS_RT_CTL